# system/systemd/mullvad-tailscale-fix.timer
[Unit]
Description=Reapply Mullvad/Tailscale coexistence rules periodically

[Timer]
OnBootSec=45s
OnUnitActiveSec=2min
OnUnitInactiveSec=2min
Persistent=true
Unit=mullvad-tailscale-fix.service

[Install]
WantedBy=timers.target
