# apps/.config/systemd/user/rfkill-guard.service
[Unit]
Description=Keep Wi-Fi unblocked on accidental rfkill toggles
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/rfkill-guard daemon
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
